# pkg/units/cockpit-bash.socket
[Unit]
Description=Bash on localhost:1112 (in cockpit context)

[Socket]
ListenStream=127.0.0.1:1112
Accept=yes

[Install]
WantedBy=sockets.target
