# pkg/units/goss.socket
[Unit]
Description=Go Server Spec stdout socket

[Socket]
ListenStream=111
Accept=yes

[Install]
WantedBy=sockets.target
